# config/etc/profile.d/geantv.sh
#!/bin/sh

export ROOTSYS=/opt/local
export G4LEVELGAMMADATA=/data/PhotonEvaporation5.3
export G4NEUTRONXSDATA=/data/G4PARTICLEXS1.1/neutron
export G4LEDATA=/data/G4EMLOW7.7
export G4SAIDXSDATA=/data/G4SAIDDATA2.0
export G4RADIOACTIVEDATA=/data/RadioactiveDecay5.3
export G4NEUTRONHPDATA=/data/G4NDL4.5
export G4ENSDFSTATEDATA=/data/G4ENSDFSTATE2.2
export ROOT_INCLUDE_PATH=/opt/local/inc/Geant
